# scripts/claude-agent.service
[Unit]
Description=Claude Code AI Agent (tmux + Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=rankgnar
Environment=HOME=/home/rankgnar
ExecStart=/home/rankgnar/agent-homebase/scripts/agent-start.sh
ExecStop=/home/rankgnar/agent-homebase/scripts/agent-stop.sh

# Restart if the tmux session dies
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
